# home/bootstrap/crontab
0 */3 * * * rm ${HOME}/.cache/bw/session & unset BW_SESSION
