# home/bootstrap/crontab
0 */3 * * * if [ -e ${HOME}/.cache/bw/session ]; then rm ${HOME}/.cache/bw/session & unset BW_SESSION; fi
